# BUILDROOT/distro-release-4.2-0.4.x86_64/etc/profile.d/10distro-release.sh
if [ -r /etc/sysconfig/system ]; then
    . /etc/sysconfig/system
    export META_CLASS
else
    export META_CLASS=unknown
fi
