# BUILDROOT/distro-release-4.2-0.4.x86_64/etc/profile.d/10distro-release.csh
if ( -r /etc/sysconfig/system ) then
    eval set SECURITY=3;
set CLASS=beginner;
set LIBSAFE=no;
set META_CLASS=download;
    setenv META_CLASS download
else
    setenv META_CLASS unknown
endif
